# .config/systemd/user/bg.service
[Unit]
Description=Change the background

[Service]
Type=oneshot
ExecStart=/home/ju/.bin/batground.sh wave1.jpg
