# .config/systemd/user/bg.timer
[Unit]
Description=Change the background

[Timer]
OnUnitActiveSec=10s
OnBootSec=10s
Unit=bg.service

[Install]
WantedBy=timers.target
